# src/FXArch/systemd/passgit.service
[Unit]
Description=Git automation for Password-Store

[Service]
ExecStart=/home/zpx/.local/bin/passgit
